# systemd/akatsukisoap-qcluster.service
[Unit]
Description=AkatsukiSoap Django-Q2 Worker
After=network.target akatsukisoap.service

[Service]
User=www-data
Group=www-data
WorkingDirectory=/opt/akatsukisoap
EnvironmentFile=/opt/akatsukisoap/.env
ExecStart=/opt/akatsukisoap/venv/bin/python manage.py qcluster
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
